# root/etc/init.d/kodexplorer
#!/bin/sh /etc/rc.common
# Copyright (C) 2018-2020 Lienol <lawlienol@gmail.com>

START=99
USE_PROCD=1  # 声明使用 procd 模式

CONFIG="kodexplorer"
TEMP_PATH="/var/etc/kodexplorer"
NGINX_CONFIG="$TEMP_PATH/nginx.conf"
PHP_FPM_CONFIG="$TEMP_PATH/php-fpm.conf"
PHP_CONFIG="/etc/php.ini"
PHP_BACKUP_CONFIG="/etc/php.ini.backup"

config_t_get() {
	local index=0
	[ -n "$4" ] && index=$4
	local ret=$(uci get $CONFIG.@$1[$index].$2 2>/dev/null)
	echo ${ret:=$3}
}

# 保留原有的配置生成函数 [cite: 1, 11]
gen_nginx_config() {
	port=$(config_t_get global port)
	https_port=$(config_t_get global https_port)
	project_directory=$(config_t_get global project_directory)
	ssl_certificate=$(config_t_get global ssl_certificate)
	ssl_certificate_key=$(config_t_get global ssl_certificate_key)
	https_enable=$(config_t_get global https_enable)

	cat <<-EOF > $1
		user root root;
		worker_processes  1;
		pid /var/run/nginx_kodexplorer.pid;
		events {
		    worker_connections  1024;
		}
		http {
		    include /etc/nginx/mime.types;
		    sendfile on;
		    keepalive_timeout 65;
		    server {
		        listen $port;
		        listen [::]:$port;
		        server_name localhost;
	EOF

	if [ "$https_enable" = "1" ]; then
		cat <<-EOF >> $1
		        listen $https_port ssl;
		        listen [::]:$https_port ssl;
		        ssl_certificate $ssl_certificate;
		        ssl_certificate_key $ssl_certificate_key;
		        ssl_session_cache shared:SSL:1m;
		        ssl_session_timeout 10m;
		        ssl_ciphers HIGH:!aNULL:!MD5;
		        ssl_prefer_server_ciphers on;
		EOF
	fi

	cat <<-EOF >> $1
		        location / {
		            root $project_directory;
		            index index.html index.htm index.php;
		        }

		        error_page 500 502 503 504 /50x.html;
		        location = /50x.html {
		            root html;
		        }

		        location ~ \.php(.*)$ {
		            root $project_directory;
		            try_files \$uri =404;
		            fastcgi_pass unix:/var/run/php8-fpm.sock;
		            fastcgi_index index.php;
		            fastcgi_split_path_info ^(.+\.php)(.*)$;
		            fastcgi_param SCRIPT_FILENAME \$document_root\$fastcgi_script_name;
		            fastcgi_param PATH_INFO \$fastcgi_path_info;
		            include /etc/nginx/fastcgi_params;
		        }
		    }
		}
	EOF
}

gen_php_config() {
	storage_device_path=$(config_t_get global storage_device_path)
	memory_limit=$(config_t_get global memory_limit)
	post_max_size=$(config_t_get global post_max_size)
	upload_max_filesize=$(config_t_get global upload_max_filesize)
	
    [ -f "$PHP_CONFIG" ] && cp $PHP_CONFIG $PHP_BACKUP_CONFIG

	cat <<-EOF >$PHP_CONFIG
		[PHP]
		engine = On
		precision = 12
		output_buffering = Off
		zlib.output_compression = Off
		implicit_flush = Off
		open_basedir = $storage_device_path:/tmp/:/proc/:/usr/bin/
		disable_functions =
		memory_limit = $memory_limit
		post_max_size = $post_max_size
		extension_dir = "/usr/lib/php8"
		cgi.fix_pathinfo=1
		file_uploads = On
		upload_max_filesize = $upload_max_filesize
		max_file_uploads = 20
		allow_url_fopen = On
		default_socket_timeout = 60
	EOF

	cat <<-EOF >$PHP_FPM_CONFIG
		[global]
		pid = /var/run/kodexplorer_php8-fpm.pid
		error_log = /var/log/kodexplorer_php8-fpm.log
		[www]
		user = root
		listen = /var/run/php8-fpm.sock
		listen.mode = 0666
		pm = dynamic
		pm.max_children = 5
		pm.start_servers = 2
		pm.min_spare_servers = 1
		pm.max_spare_servers = 3
		chdir = /
	EOF
}

# Procd 服务启动函数
start_service() {
	ENABLED=$(config_t_get global enable 0)
	[ "$ENABLED" = "0" ] && return 0

	mkdir -p $TEMP_PATH /var/log/nginx /var/lib/nginx
	
	# 生成配置
	gen_php_config
	gen_nginx_config $NGINX_CONFIG

	# 启动 PHP8-FPM 实例
	procd_open_instance "php-fpm"
	procd_set_param command /usr/bin/php8-fpm -R -y $PHP_FPM_CONFIG -g "/var/run/php8-fpm.pid" --nodaemonize
	procd_set_param respawn # 进程崩溃自动重启
	procd_set_param file $PHP_FPM_CONFIG # 监控配置文件变化
	procd_set_param stdout 1 # 将日志输出到系统 log
	procd_set_param stderr 1
	procd_close_instance

	# 启动 Nginx 实例
	procd_open_instance "nginx"
	# 注意：必须添加 -g "daemon off;" 
	procd_set_param command /usr/sbin/nginx -c $NGINX_CONFIG -g "daemon off;"
	procd_set_param respawn
	procd_set_param file $NGINX_CONFIG
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

# 服务关闭时的清理动作
service_stopped() {
    # 恢复 PHP 原始配置 
	[ -f "$PHP_BACKUP_CONFIG" ] && {
		mv $PHP_BACKUP_CONFIG $PHP_CONFIG
	}
    rm -f /var/run/php8-fpm.sock
}

# 注册服务触发器，当 uci 配置改变时自动重启
service_triggers() {
	procd_add_reload_trigger "$CONFIG"
}